# conf/dynette-regen-named-conf.path
[Path]
Unit=__APP__-regen-named-conf.service
PathModified=__DATA_DIR__/domains.sql.flag

[Install]
WantedBy=multi-user.target
