# conf/dynette-regen-named-conf.path
[Path]
Unit=__APP__-regen-named-conf.service
PathChanged=__DATA_DIR__/db

[Install]
WantedBy=multi-user.target
